# install_service/ephemera.service
[Unit]
SourcePath=/etc/init.d/ephemera
[Service]
ExecStart=/etc/init.d/ephemera start
ExecStop=/etc/init.d/ephemera stop